# deploy/cursorbot-b4-5m.service
[Unit]
Description=Cursorbot B4 Spread Runner (5-min BTC, Polymarket)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/cursorbot
ExecStart=/usr/bin/node dist/b4-5m/spread-run.js
Restart=always
RestartSec=10
EnvironmentFile=/root/cursorbot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cursorbot-b4

[Install]
WantedBy=multi-user.target
