# deploy/cursorbot-paper.service
[Unit]
Description=Cursorbot B4/B5 Paper Trader (no real orders)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/cursorbot
EnvironmentFile=/root/cursorbot/.env
ExecStart=/usr/bin/node dist/b4-5m/paper-run.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cursorbot-paper

[Install]
WantedBy=multi-user.target
