# deploy/cursorbot-b5.service
[Unit]
Description=Cursorbot B5 Basket Bot (5m/15m BTC/ETH, D3 only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/cursorbot
ExecStart=/usr/bin/node dist/b5-basket/run.js
Restart=always
RestartSec=10
EnvironmentFile=/root/cursorbot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cursorbot-b5

[Install]
WantedBy=multi-user.target
