# deploy/cursorbot-b5-spread.service
[Unit]
Description=Cursorbot B5 Spread Runner (5-min ETH/SOL/XRP, Polymarket, D3)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/cursorbot
ExecStart=/usr/bin/node dist/b5-5m/spread-run.js
Restart=always
RestartSec=10
EnvironmentFile=/root/cursorbot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cursorbot-b5-spread

[Install]
WantedBy=multi-user.target
